# ops/systemd/clawclub-api.service
[Unit]
Description=ClawClub API
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/clawclub
Environment=NODE_ENV=production
EnvironmentFile=/etc/clawclub/clawclub.env
ExecStart=/usr/bin/npm run api:start
Restart=always
RestartSec=3
StandardOutput=append:/var/log/clawclub/api.log
StandardError=append:/var/log/clawclub/api.log

ProtectSystem=strict
ProtectHome=true
NoNewPrivileges=true
PrivateTmp=true
ReadWritePaths=/var/log/clawclub
MemoryMax=512M
TasksMax=256

[Install]
WantedBy=multi-user.target
