# ops/systemd/clawclub-worker.service
[Unit]
Description=ClawClub delivery worker
After=network-online.target clawclub-api.service
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/clawclub
EnvironmentFile=/etc/clawclub/clawclub.env
ExecStart=/bin/bash -lc 'while true; do npm run api:worker -- --worker-key hetzner-main --max-runs 10; sleep 5; done'
Restart=always
RestartSec=3
StandardOutput=append:/var/log/clawclub/worker.log
StandardError=append:/var/log/clawclub/worker.log

[Install]
WantedBy=multi-user.target
